# crontab
* * * * * root /opt/certbot/venv/bin/certbot --standalone certonly --agree-tos --email budde377@gmail.com -d laundree.io >> /var/log/cron.log 2>&1
